# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
priority=1
command=/usr/sbin/uwsgi --ini /app/uwsgi.ini --die-on-term --need-app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:nginx]
priority=10
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
# Graceful stop, see http://nginx.org/en/docs/control.html
stopsignal=QUIT
